# launchagent/com.jarvis.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/daijiromatsuokam1/.bun/bin/bun</string>
        <string>--env-file=.env</string>
        <string>run</string>
        <string>src/index.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/daijiromatsuokam1/claude-telegram-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/daijiromatsuokam1/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/daijiromatsuokam1</string>
    </dict>

    <!-- ログイン時に自動起動 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- プロセス死亡時に自動再起動 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 再起動間隔の最小秒数（クラッシュループ防止） -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- ログ出力 -->
    <key>StandardOutPath</key>
    <string>/Users/daijiromatsuokam1/claude-telegram-bot/logs/bot-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/daijiromatsuokam1/claude-telegram-bot/logs/bot-launchd.err</string>
</dict>
</plist>
